# build/stage3/03-autostart/files/alsa-multi-card.service
[Unit]
Description=Generate Multi-Card ALSA Configuration
Documentation=https://github.com/benpaddlejones/Raspberry_HDMI_Tester
After=sound.target
Before=hdmi-audio-ready.service

[Service]
Type=oneshot
ExecStart=/usr/local/bin/generate-multi-card-asound
RemainAfterExit=yes
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=sound.target
