# build/stage3/03-autostart/files/audio-test.service
[Unit]
Description=HDMI Audio Test - Stereo and 5.1 Surround FLAC Loop
After=local-fs.target sound.target
Wants=sound.target

[Service]
Type=simple
User=pi
SupplementaryGroups=audio video render

# Pre-flight checks
ExecStartPre=/bin/bash -c 'echo "=== Audio Test Service Starting ===" && \
    echo "Timestamp: $(date)" && \
    echo "Checking test audio files..." && \
    if [ ! -f /opt/hdmi-tester/stereo.flac ]; then \
        echo "ERROR: Stereo FLAC not found at /opt/hdmi-tester/stereo.flac" && \
        exit 1; \
    fi && \
    if [ ! -f /opt/hdmi-tester/surround51.flac ]; then \
        echo "ERROR: 5.1 Surround FLAC not found at /opt/hdmi-tester/surround51.flac" && \
        exit 1; \
    fi && \
    echo "✓ Stereo FLAC exists ($(stat -c%%s /opt/hdmi-tester/stereo.flac) bytes)" && \
    echo "✓ 5.1 Surround FLAC exists ($(stat -c%%s /opt/hdmi-tester/surround51.flac) bytes)" && \
    echo "Launching audio-test..."'

# Run audio-test command (loop both FLAC files) using VLC
# Use AUDIODEV environment variable to specify ALSA device for HDMI audio
ExecStart=/bin/bash -c 'HDMI_DEVICE=$(/opt/hdmi-tester/detect-hdmi-audio) && \
    export AUDIODEV=${HDMI_DEVICE} && \
    cvlc --loop --no-video --aout=alsa \
    --quiet --no-audio-time-stretch --audio-replay-gain-mode=none \
    /opt/hdmi-tester/stereo.flac /opt/hdmi-tester/surround51.flac'

# Log exit status
ExecStopPost=/bin/bash -c 'echo "=== Audio Test Service Stopped ===" && \
    echo "Exit code: $EXIT_CODE" && \
    if [ "$SERVICE_RESULT" != "success" ]; then \
        echo "ERROR: Service failed - check logs"; \
    fi'

Restart=always
RestartSec=5
StandardOutput=journal
StandardError=journal
KillMode=mixed
TimeoutStopSec=10

[Install]
WantedBy=multi-user.target
