# build/stage3/03-autostart/files/audio-test.service
[Unit]
Description=HDMI Audio Test - Stereo and 5.1 Surround FLAC Loop
After=local-fs.target sound.target systemd-user-sessions.service hdmi-audio-ready.service
Wants=sound.target
Requires=hdmi-audio-ready.service

# CRITICAL: Wait for ALSA/bcm2835 driver initialization to complete
# The bcm2835 audio driver needs time to communicate with VideoCore firmware
# Starting too early causes VCHI timeouts and kernel crashes
Requires=systemd-user-sessions.service

[Service]
Type=simple
User=pi
SupplementaryGroups=audio video render

# Pre-flight checks
ExecStartPre=/bin/bash -c 'echo "=== Audio Test Service Starting ===" && \
    echo "Timestamp: $(date)" && \
    echo "Waiting 3 seconds for ALSA driver initialization..." && \
    sleep 3 && \
    echo "Checking test audio files..." && \
    if [ ! -f /opt/hdmi-tester/stereo.flac ]; then \
        echo "ERROR: Stereo FLAC not found at /opt/hdmi-tester/stereo.flac" && \
        exit 1; \
    fi && \
    if [ ! -f /opt/hdmi-tester/surround51.flac ]; then \
        echo "ERROR: 5.1 Surround FLAC not found at /opt/hdmi-tester/surround51.flac" && \
        exit 1; \
    fi && \
    echo "✓ Stereo FLAC exists ($(stat -c%%s /opt/hdmi-tester/stereo.flac) bytes)" && \
    echo "✓ 5.1 Surround FLAC exists ($(stat -c%%s /opt/hdmi-tester/surround51.flac) bytes)" && \
    echo "Launching audio-test..."'

# Run the dedicated audio-test wrapper (handles device detection and logging)
ExecStart=/opt/hdmi-tester/audio-test

# Log exit status
ExecStopPost=/bin/bash -c 'echo "=== Audio Test Service Stopped ===" && \
    echo "Exit code: $EXIT_CODE" && \
    if [ "$SERVICE_RESULT" != "success" ]; then \
        echo "ERROR: Service failed - check logs"; \
    fi'

Restart=always
RestartSec=5
StandardOutput=journal
StandardError=journal
KillMode=mixed
TimeoutStopSec=10

[Install]
WantedBy=multi-user.target
