# build/stage3/03-autostart/files/full-test.service
[Unit]
Description=HDMI Full Test - Both Videos Loop
After=local-fs.target sound.target
Wants=sound.target

[Service]
Type=simple
User=pi
Group=audio

# Pre-flight checks
ExecStartPre=/bin/bash -c 'echo "=== Full Test Service Starting ===" && \
    echo "Timestamp: $(date)" && \
    echo "Checking test video files..." && \
    if [ ! -f /opt/hdmi-tester/image-test.webm ]; then \
        echo "ERROR: Image test video not found" && \
        exit 1; \
    fi && \
    if [ ! -f /opt/hdmi-tester/color_test.webm ]; then \
        echo "ERROR: Color test video not found" && \
        exit 1; \
    fi && \
    echo "✓ Both test videos exist" && \
    echo "Launching test-both-loop..."'

# Run test-both-loop command (both videos in sequence, infinite loop)
ExecStart=/bin/bash -c 'while true; do \
    echo "Playing: image-test.webm (optimized resolution)..." && \
    mpv --no-osd-bar --osd-level=0 \
        --hwdec=no --vo=drm --drm-connector=HDMI-A-1 \
        --cache=yes --demuxer-max-bytes=100M --no-config \
        /opt/hdmi-tester/image-test.webm 2>/dev/null && \
    echo "Playing: color_test.webm (fullscreen stretched)..." && \
    mpv --fs --no-keepaspect --no-osd-bar --osd-level=0 \
        --hwdec=no --vo=drm --drm-connector=HDMI-A-1 \
        --cache=yes --demuxer-max-bytes=100M --no-config \
        /opt/hdmi-tester/color_test.webm 2>/dev/null && \
    echo "Restarting loop..."; \
done'

# Log exit status
ExecStopPost=/bin/bash -c 'echo "=== Full Test Service Stopped ===" && \
    echo "Exit code: $EXIT_CODE" && \
    if [ "$SERVICE_RESULT" != "success" ]; then \
        echo "ERROR: Service failed - check logs"; \
    fi'

Restart=always
RestartSec=5
StandardOutput=journal
StandardError=journal
KillMode=mixed
TimeoutStopSec=10

[Install]
WantedBy=multi-user.target
