# build/stage3/03-autostart/files/hd-audio-test.service
[Unit]
Description=HDMI HD Audio Test - Image Loop
After=local-fs.target sound.target
Wants=sound.target

[Service]
Type=simple
User=pi
Group=audio

# Pre-flight checks
ExecStartPre=/bin/bash -c 'echo "=== HD Audio Test Service Starting ===" && \
    echo "Timestamp: $(date)" && \
    echo "Checking test video file..." && \
    if [ ! -f /opt/hdmi-tester/image-test.webm ]; then \
        echo "ERROR: Test video not found at /opt/hdmi-tester/image-test.webm" && \
        exit 1; \
    fi && \
    echo "✓ Test video exists ($(stat -c%%s /opt/hdmi-tester/image-test.webm) bytes)" && \
    echo "Launching test-image-loop..."'

# Run test-image-loop command (optimized resolution with audio)
ExecStart=/usr/bin/mpv --loop=inf --no-osd-bar --osd-level=0 --no-input-default-bindings \
    --hwdec=no --vo=drm --drm-connector=HDMI-A-1 \
    --cache=yes --demuxer-max-bytes=100M --no-config \
    /opt/hdmi-tester/image-test.webm

# Log exit status
ExecStopPost=/bin/bash -c 'echo "=== HD Audio Test Service Stopped ===" && \
    echo "Exit code: $EXIT_CODE" && \
    if [ "$SERVICE_RESULT" != "success" ]; then \
        echo "ERROR: Service failed - check logs"; \
    fi'

Restart=always
RestartSec=5
StandardOutput=journal
StandardError=journal
KillMode=mixed
TimeoutStopSec=10

[Install]
WantedBy=multi-user.target
