# build/stage3/03-autostart/files/hdmi-audio.service
[Unit]
Description=HDMI Audio Test - ALSA Output
After=sound.target

[Service]
Type=simple
User=pi
Group=audio
# Use ALSA for direct audio output to HDMI
ExecStart=/usr/bin/mpv --loop=inf --no-video --ao=alsa --audio-device=alsa/plughw:CARD=vc4hdmi,DEV=0 --volume=100 /opt/hdmi-tester/audio.mp3
Restart=always
RestartSec=15
StandardOutput=journal
StandardError=journal
KillMode=mixed
TimeoutStopSec=10

[Install]
WantedBy=multi-user.target
