# build/stage3/03-autostart/files/hdmi-audio-config.service
[Unit]
Description=Generate dynamic ALSA config for HDMI audio
Requires=sound.target
After=sound.target

[Service]
Type=oneshot
ExecStartPre=/bin/sleep 2
ExecStart=/usr/local/bin/generate-asound-conf
RemainAfterExit=yes
Restart=on-failure
RestartSec=5

[Install]
WantedBy=multi-user.target
