# build/stage3/03-autostart/files/hdmi-display.service
[Unit]
Description=HDMI Test Pattern Display (Framebuffer)
After=local-fs.target

[Service]
Type=simple
User=root
# Use fbi with -T 1 for HDMI output, -a for auto-zoom, --noverbose to suppress output, -d to keep displayed
ExecStart=/usr/bin/fbi -T 1 -a --noverbose -d /opt/hdmi-tester/image.png
Restart=always
RestartSec=10
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target
