# build/stage3/03-autostart/files/pixel-audio-test-vlc.service
[Unit]
Description=HDMI Pixel Audio Test - Color Fullscreen (VLC)
After=local-fs.target sound.target
Wants=sound.target

[Service]
Type=simple
User=pi
Group=audio

# Pre-flight checks
ExecStartPre=/bin/bash -c 'echo "=== Pixel Audio Test Service (VLC) Starting ===" && \
    echo "Timestamp: $(date)" && \
    echo "Checking test video file..." && \
    if [ ! -f /opt/hdmi-tester/color-test.webm ]; then \
        echo "ERROR: Test video not found at /opt/hdmi-tester/color-test.webm" && \
        exit 1; \
    fi && \
    echo "✓ Test video exists ($(stat -c%%s /opt/hdmi-tester/color-test.webm) bytes)" && \
    echo "Launching test-color-fullscreen with VLC..."'

# Run test-color-fullscreen command (fullscreen stretched with audio) using VLC
# --no-autoscale forces stretching to fill screen (for pixel testing)
ExecStart=/usr/bin/cvlc --loop --fullscreen --no-autoscale \
    --no-video-title-show --no-osd --quiet \
    --vout=drm --alsa-audio-device hw:CARD=vc4hdmi,DEV=0 \
    --no-audio-time-stretch --audio-replay-gain-mode=none \
    /opt/hdmi-tester/color-test.webm

# Log exit status
ExecStopPost=/bin/bash -c 'echo "=== Pixel Audio Test Service (VLC) Stopped ===" && \
    echo "Exit code: $EXIT_CODE" && \
    if [ "$SERVICE_RESULT" != "success" ]; then \
        echo "ERROR: Service failed - check logs"; \
    fi'

Restart=always
RestartSec=5
StandardOutput=journal
StandardError=journal
KillMode=mixed
TimeoutStopSec=10

[Install]
WantedBy=multi-user.target
