# build/stage3/03-autostart/files/pixel-test.service
[Unit]
Description=HDMI Pixel Test - Color Fullscreen
After=local-fs.target sound.target
Wants=sound.target

[Service]
Type=simple
User=pi
SupplementaryGroups=audio video render

# Pre-flight checks
ExecStartPre=/bin/bash -c 'echo "=== Pixel Test Service Starting ===" && \
    echo "Timestamp: $(date)" && \
    echo "Checking test video files..." && \
    if [ ! -f /opt/hdmi-tester/color-test.webm ] && [ ! -f /opt/hdmi-tester/color-test.mp4 ]; then \
        echo "ERROR: No test video found (checked .webm and .mp4)" && \
        exit 1; \
    fi && \
    echo "✓ Test video files available" && \
    echo "Launching pixel-test script (auto-detects format)..."'

# Run pixel-test script (auto-detects Pi model and selects appropriate format)
ExecStart=/opt/hdmi-tester/pixel-test

# Log exit status
ExecStopPost=/bin/bash -c 'echo "=== Pixel Test Service Stopped ===" && \
    echo "Exit code: $EXIT_CODE" && \
    if [ "$SERVICE_RESULT" != "success" ]; then \
        echo "ERROR: Service failed - check logs"; \
    fi'

Restart=always
RestartSec=5
StandardOutput=journal
StandardError=journal
KillMode=mixed
TimeoutStopSec=10

[Install]
WantedBy=multi-user.target
