# build/stage3/03-autostart/files/test-notvideo.service
[Unit]
Description=HDMI Static Image + Audio Test Service
Documentation=https://github.com/benpaddlejones/Raspberry_HDMI_Tester
After=multi-user.target
Conflicts=hdmi-test.service pixel-test.service full-test.service image-test.service

[Service]
Type=simple
ExecStart=/opt/hdmi-tester/test-notvideo
Restart=always
RestartSec=5
StandardOutput=journal
StandardError=journal
SyslogIdentifier=test-notvideo

# Resource limits
CPUQuota=80%
MemoryMax=512M

# Security
User=pi
SupplementaryGroups=audio video render
NoNewPrivileges=true
PrivateTmp=true

[Install]
WantedBy=multi-user.target
