# build/stage3/04-boot-config/files/audio-health.service
[Unit]
Description=Audio Health Check (Diagnostic)
After=sound.target systemd-user-sessions.service
ConditionPathExists=/opt/hdmi-tester/check-audio-health

[Service]
Type=oneshot
ExecStart=/opt/hdmi-tester/check-audio-health
RemainAfterExit=yes
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target
