# build/stage3/04-boot-config/files/fix-cmdline.timer
[Unit]
Description=Deferred first-boot cmdline cleanup
Documentation=https://github.com/benpaddlejones/Raspberry_HDMI_Tester

[Timer]
OnBootSec=1min
AccuracySec=5s
Unit=fix-cmdline.service
Persistent=no

[Install]
WantedBy=timers.target
